# sparse/usr/lib/systemd/system/yamuisplash.service.d/50-vendor-hwcomposer.conf
[Service]
ExecStartPre=-/system/bin/stop vendor.qti.hardware.display.composer

# stop first in case something else managed to start it
ExecStopPost=-/system/bin/stop vendor.qti.hardware.display.composer
ExecStopPost=-/system/bin/start vendor.qti.hardware.display.composer
