# sparse/usr/lib/systemd/system/sailfish-unlock-agent.service.d/50-vendor-hwcomposer.conf
[Service]
# stop hwcomposer before unlock ui
ExecStartPre=-/system/bin/stop vendor.qti.hardware.display.composer

# start hwcomposer after unlock ui, but not on upgrade target
ExecStart=-/bin/sh -c '/usr/bin/test -f /tmp/os-update-running || /system/bin/start vendor.qti.hardware.display.composer'
