# sparse/usr/lib/systemd/user/lipstick.service.d/50-vendor-hwcomposer.conf
[Service]
# make unlock ui exit
ExecStartPre=/usr/sbin/dummy_compositor
ExecStopPost=/usr/bin/setprop ctl.restart vendor.qti.hardware.display.composer
